# Storj3Monitor/Storj3Monitor.service
[Unit]
Description=Storj v3 monitor by Krey
Requires=network.target
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/pwsh /etc/scripts/storj/Storj3Monitor.ps1 -c /etc/scripts/storj/Storj3Monitor.my.conf monitor
ExecStop=/bin/kill --signal SIGINT ${MAINPID}

[Install]
WantedBy=multi-user.target
